# base/s6/services.d/chisel/run
#!/usr/bin/with-contenv bash

chisel server --port 8888 --proxy http://localhost