# s6/services.d/instructions/run
#!/usr/bin/with-contenv bash

cd /instructions

python3 -m http.server 80